# systemd/kiosk-update.service
[Unit]
Description=Hive Kiosk self-update from origin/main

[Service]
Type=oneshot
ExecStart=/opt/weather_kiosk/bin/kiosk-update.sh
